# templates/default/nodejs.upstart.conf.erb
# This is an upstart (http://upstart.ubuntu.com/) script
# to run the node.js server on system boot and make it
# manageable with commands such as
# 'start node-app' and 'stop node-app'
#
# This script is to be placed in /etc/init to work with upstart.
#
# Internally the 'initctl' command is used to manage:
# initctl help
# initctl status node-app
# initctl reload node-app
# initctl start node-app

description "node.js forever server for node-app"
version "1.0"

expect fork

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:

start on filesystem and started networking
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

setuid <%= @user %>
setgid <%= @user %>

env HOME=<%= @app_dir %>
<% @env.each_pair do | variable, value | %>
env <%= variable %>=<%= value %>
<% end %>

script
    export HOME=$HOME
  <% @env.each_pair do | variable, value | %>
    export <%= variable %>=$<%= variable %>
  <% end %>
    echo $PORT
    chdir $HOME
    exec <%= @binary_path %> <%= @entry %> > <%= @app_dir %>/logs/node.log &
end script

#post-start script
#   # Optionally put a script here that will notifiy you node has (re)started
#   # /root/bin/hoptoad.sh "node.js has started!"
#end script
