# templates/default/nodejs.initd.erb
#!/bin/sh

#
# chkconfig: 35 99 99
# description: Node.js /home/nodejs/sample/app.js
#

. /etc/rc.d/init.d/functions

USER="<%= @user %>"
GROUP="<%= @group %>"

BINARY="<%= @binary_path %>"

ROOT_DIR="<%= @app_dir %>"

SERVER="$ROOT_DIR/<%= @entry %>"
LOG_FILE="<%= @app_dir %>/logs/node.log"

LOCK_FILE="/var/lock/subsys/<%= @app_name %>"

do_start()
{
        export HOME=$ROOT_DIR
        if [ ! -f "$LOCK_FILE" ] ; then
                echo -n $"Starting $SERVER: "
                daemon --user "$USER" "$BINARY $SERVER >> $LOG_FILE &" && echo_success || echo_failure
                RETVAL=$?
                echo
                [ $RETVAL -eq 0 ] && touch $LOCK_FILE
        else
                echo "$SERVER is locked. Remove $LOCK_FILE if you think this is wrong."
                RETVAL=1
        fi
}
do_stop()
{
        echo -n $"Stopping $SERVER: "
        pid=`ps -aefw | grep "$DAEMON $SERVER" | grep -v " grep " | awk '{print $2}'`
        pkill -TERM -P $pid > /dev/null 2>&1 && echo_success || echo_failure
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && rm -f $LOCK_FILE
}

case "$1" in
        start)
                do_start
                ;;
        stop)
                do_stop
                ;;
        restart)
                do_stop
                do_start
                ;;
        status)
                status $BINARY
                ;;
        *)
                echo "Usage: $0 {start|stop|restart}"
                RETVAL=1
esac

exit $RETVAL
