# templates/default/nodejs.upstart.conf.erb
#!upstart
description "Node.js Application Server"

start on startup
stop on shutdown

expect fork

<% unless @user.nil? -%>
setuid <%= @user %>
<% end -%>
<% unless @group.nil? -%>
setgid <%= @group %>
<% end -%>
chdir <%= @app_dir %>

script
  cd <%= @app_dir %>
  exec forever start \
    --pidFile <%= @app_dir %>/forever.pid \
    --minUptime 5000 \
    --spinSleepTime 2000 \
    -a \
    -l <%= @app_dir %>/<%= @app_name %>.log \
    <%= @app_dir %>/<%= @entry %>
end script

pre-stop script
  exec forever stop <%= @app_dir %>/<%= @entry %>
end script